# lib/kafka/09-topic-config.conf
[program:topic-config]
user=nobody
command=bash -c 'eval /usr/local/share/landoop/wait-scripts/wait-for-kafka.sh; exec /usr/libexec/pccx/kafka-topic-config'
redirect_stderr=true
stdout_logfile=/var/log/topic-config.log
startretries=1
